# deploy/systemd/easytech-meta-oauth.service
[Unit]
Description=EasyMarketingOne Meta OAuth (Facebook + Instagram)
After=network.target

[Service]
Type=simple
User=papichulo
WorkingDirectory=/opt/easytech_marketing
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/easytech_marketing/venv/bin/python3 /opt/easytech_marketing/Motor_Tecnico/meta_oauth.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
