# deploy/systemd/easytech-accio-engine-test.service
[Unit]
Description=EasyMarketingOne Accio Engine TEST
After=network.target

[Service]
Type=simple
User=papichulo
WorkingDirectory=/opt/easytech_marketing
Environment=PATH=/opt/easytech_marketing/venv/bin
Environment=PYTHONUNBUFFERED=1
Environment=ACCIO_ENV=test
EnvironmentFile=-/opt/easytech_marketing/deploy/secrets/.env.test
ExecStart=/opt/easytech_marketing/venv/bin/python3 /opt/easytech_marketing/Motor_Tecnico/accio_engine/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
